# orchestration/error-recovery.cron
*/15 * * * * cd /path/to/workspace-master && node orchestration/error-recovery.mjs --auto-recover
