# orchestration/health-check.cron
*/5 * * * * cd /path/to/workspace-master && node orchestration/health-check.mjs --alert-on-error
